# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:pathway]
command=python /app/main.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/pathway.err.log
stdout_logfile=/var/log/pathway.out.log
environment=PYTHONUNBUFFERED=1

[program:postgres-writer]
command=python /app/postgres_writer.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/postgres-writer.err.log
stdout_logfile=/var/log/postgres-writer.out.log
environment=PYTHONUNBUFFERED=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx/nginx.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log